# projects/ai-team/deploy/com.ailab.yewon_selfheal.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key><string>com.ailab.yewon_selfheal</string>
    <!-- 매일 08:00 자가 점검→안전 자동복구→결과만 텔레그램 보고 -->
    <key>ProgramArguments</key>
    <array>
        <string>__PYTHON__</string>
        <string>/Users/junholee/ai_lab/projects/ai-team/skills/예원_CEO/tools/yewon_self_heal.py</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key><integer>8</integer>
        <key>Minute</key><integer>0</integer>
    </dict>
    <key>RunAtLoad</key><false/>
    <key>WorkingDirectory</key><string>/Users/junholee/ai_lab</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONUTF8</key><string>1</string>
    </dict>
    <key>StandardOutPath</key><string>/Users/junholee/ai_lab/output/bot_logs/yewon_selfheal.out.log</string>
    <key>StandardErrorPath</key><string>/Users/junholee/ai_lab/output/bot_logs/yewon_selfheal.err.log</string>
</dict>
</plist>
