# projects/ai-team/scripts/launchd/com.ailab.somi.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ailab.somi</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>-u</string>
        <string>/Users/junholee/ai_lab/projects/ai-team/skills/소미_메일매니저/tools/gmail_manager.py</string>
        <string>--max</string>
        <string>100</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/junholee/ai_lab</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONUTF8</key>
        <string>1</string>
        <key>HOME</key>
        <string>/Users/junholee</string>
    </dict>
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Hour</key><integer>8</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>14</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>21</integer><key>Minute</key><integer>0</integer></dict>
    </array>
    <key>StandardOutPath</key>
    <string>/Users/junholee/ai_lab/output/bot_logs/somi.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/junholee/ai_lab/output/bot_logs/somi.err.log</string>
</dict>
</plist>
