# projects/ai-team/scripts/launchd/com.ailab.youngsuk.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ailab.youngsuk</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/junholee/ai_lab/projects/ai-team/skills/영숙_비서/tools/youngsuk_launcher.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/junholee/ai_lab</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONUTF8</key>
        <string>1</string>
        <key>HOME</key>
        <string>/Users/junholee</string>
    </dict>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>30</integer>
    <key>StandardOutPath</key>
    <string>/Users/junholee/ai_lab/output/trading_logs/youngsuk_daemon.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/junholee/ai_lab/output/trading_logs/youngsuk_daemon.err.log</string>
</dict>
</plist>
